# addon/planthub/rootfs/etc/services.d/planthub/run
#!/usr/bin/with-contenv sh
set -e

echo "[PlantHub] Starting PlantHub..."

cd /opt/planthub
exec ./PlantHub.Web